# app/scrapping_data/app/crontab
# Crontab pour exécuter le script de scraping tous les jours à 21h00
0 21 * * * python /app/scrapping.py >> /var/log/cron.log 2>&1

#